# scripts/sec26-button-dispatch.service
[Unit]
Description=SEC26 Button → Gitea workflow_dispatch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ieee
WorkingDirectory=/home/ieee/SEC26
Environment=PYTHONUNBUFFERED=1
# Configure these via Environment or edit values in the script's defaults
Environment=GITEA_URL=https://gitea.syndric.org
Environment=GITEA_OWNER=ieeeucf
Environment=GITEA_REPO=SEC26
Environment=GITEA_WORKFLOW=pi-deploy.yml
Environment=GITEA_REF=prod
Environment=GITEA_TOKEN=REPLACE_WITH_TOKEN
Environment=BTN_PIN=22
ExecStart=/usr/bin/python3 /home/ieee/SEC26/scripts/button-trigger-workflow.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
